# dockerfiles/etc/supervisor/conf.d/gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn --workers=2 main:app -b 0.0.0.0:8002
directory=/home/admin/web/pureblog
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
